# deploy/inclu-ia.service
[Unit]
Description=Inclu-IA captions server
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
User=__RUN_USER__
Group=audio
WorkingDirectory=__REPO_DIR__/software
EnvironmentFile=-__REPO_DIR__/software/.env
ExecStart=__REPO_DIR__/software/.venv/bin/python __REPO_DIR__/software/server.py
Restart=always
RestartSec=2
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
